# deployment/local.ustadai.3003.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.ustadai.3003</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-l</string>
      <string>-c</string>
      <string>cd "/Users/ihu/Documents/New project" &amp;&amp; if [ ! -f ".next/BUILD_ID" ]; then /usr/local/bin/npm run build; fi &amp;&amp; exec /usr/local/bin/npm run start -- --hostname 127.0.0.1 --port 3003</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ihu/Documents/New project</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ihu/Documents/New project/.data/logs/launchd-3003.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ihu/Documents/New project/.data/logs/launchd-3003.err.log</string>
  </dict>
</plist>
